# systemd/minecraft-resume.service
[Unit]
Description=Handle Minecraft startup after hibernation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/minecraft-startup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
